# .docker/fpm/var/spool/cron/crontabs/www-data
# see https://doc.oroinc.com/backend/cron/
*/1 * * * *  . /usr/local/share/env.sh > /dev/null && /usr/bin/php /var/www/publicadmin/cli/cron.php >> /dev/null