# .docker/fpm/var/spool/cron/crontabs/www-data
# see https://doc.oroinc.com/backend/cron/
*/1 * * * * cd /var/www && . /usr/local/share/env.sh > /dev/null && /usr/bin/php bin/console oro:cron >> /dev/null